# conf/supervisor.conf
[program:golf_uwsgi]
command=/usr/local/bin/uwsgi --ini /home/ubuntu/golf-app/conf/uwsgi.ini
stopsignal=QUIT
autostart=true
autorestart=true
user=ubuntu